# config/supervisord/supervisord.conf
[program:rsyslogd]
command = /usr/sbin/rsyslogd -n
autostart=true
autorestart=true

[program:cron]
command = /usr/sbin/cron -f -L 15
autostart=true
autorestart=true
